# teachatty_chat_app_daphne.service
[Unit]
Description=Teachatty Chat App Daphne Service

[Service]
WorkingDirectory=/var/lib/jenkins/workspace/teachatty_chat_app/
User=jenkins
Environment=DJANGO_SETTINGS_MODULE=chatSystem.settings
ExecStart=/bin/bash /var/lib/jenkins/workspace/teachatty_chat_app/daphne_en.sh
Restart=always
KillSignal=SIGQUIT
Type=notify
NotifyAccess=all

[Install]
WantedBy=multi-user.target